# _release/sync-rtc.cron
#!/bin/sh

exec /usr/bin/rtc-utils --write-if-sync
