# _release/sync-rtc.cron
#!/bin/sh

exec /usr/bin/rtc write | logger --tag cron-rtc
